# .crontab
MAILTO=kevin@paretje.be
*/15 * * * * which newsbeuter > /dev/null 2>&1 && ! $HOME/bin/3gusage && chronic newsbeuter -x reload
