# .crontab
MAILTO=kevin@paretje.be
*/2 * * * * DISPLAY=:0 $HOME/bin/batmonitor
*/15 * * * * which newsbeuter 2>&1 && ! $HOME/bin/3gusage && chronic newsbeuter -x reload
