# .crontab
*/2 * * * * DISPLAY=:0.0 /home/kevin/bin/batmonitor
