# .crontab
MAILTO=kevin@paretje.be
*/2 * * * * DISPLAY=:0.0 $HOME/bin/batmonitor
