# .crontab
MAILTO=kevin@paretje.be
*/15 * * * * which newsboat > /dev/null 2>&1 && ! $HOME/bin/3gusage && chronic newsboat -x reload
