# .crontab
*/2 * * * * DISPLAY=:0.0 $HOME/bin/batmonitor
